# seccheck-weekly.timer
[Unit]
Description=Weekly seccheck run

[Timer]
OnCalendar=weekly
AccuracySec=6h
Persistent=true

[Install]
WantedBy=seccheck-weekly.service
